# systemd/narad.service
[Unit]
Description=Narad - Nisha Platform Intelligence Brain
Documentation=https://github.com/sdachary/narad
After=network-online.target
Wants=network-online.target

[Service]
# ── User & Permissions ────────────────────────────
Type=simple
User=deepak
WorkingDirectory=/home/deepak/Work/narad

# ── Environment ────────────────────────────────────
EnvironmentFile=/home/deepak/Work/narad/.env
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"

# ── Startup ────────────────────────────────────────
ExecStartPre=/home/deepak/Work/narad/scripts/validate-env.sh
ExecStart=/usr/bin/node /home/deepak/Work/narad/src/main.js

# ── Restart Policy ────────────────────────────────
Restart=on-failure
RestartSec=30
StartLimitInterval=600
StartLimitBurst=5

# ── Resource Limits ────────────────────────────────
MemoryLimit=2G
CPUQuota=80%

# ── Logging ────────────────────────────────────────
StandardOutput=journal
StandardError=journal
SyslogIdentifier=narad

# ── Hardening ────────────────────────────────────
NoNewPrivileges=true
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/deepak/Work/narad

[Install]
WantedBy=multi-user.target
